# config/systemd/user/checkupdates.timer
[Unit]
Description=Run checkupdate 15min Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
